# docs/supervisor.conf
; ============================================================
; Umar Bakery — Supervisor Worker Configuration
; ============================================================
; USAGE: Copy this file to /etc/supervisor/conf.d/umarbakery.conf
;        on your Linux production server. Then run:
;
;        sudo supervisorctl reread
;        sudo supervisorctl update
;        sudo supervisorctl start umarbakery-worker:*
; ============================================================

[program:umarbakery-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/umarbakery/artisan queue:work database --sleep=3 --tries=3 --max-time=3600
; ^ Change 'database' to 'redis' if you have migrated to Redis
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
; ^ Change to your server's web user (e.g. nginx, ubuntu)
numprocs=2
; ^ Number of worker processes. Start with 2 for a VPS.
redirect_stderr=true
stdout_logfile=/var/www/umarbakery/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600
